# crontab
10 13 25 * * python3 /usr/local/main.py
# crontab requires an empty line at the end of the file
